# templates/profile.d/bundle.sh
# Add vendor/bundle/.bin to the PATH
export PATH={{code_dir}}/vendor/bundle/.bin:$PATH
